# scripts/mochitto-client@.service
[Unit]
Description=Mochitto Voice Assistant Client
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/mochitto
ExecStart=/home/%i/.local/bin/uv run mochitto-client
Restart=on-failure
RestartSec=5

# .env から環境変数を読み込み
EnvironmentFile=/home/%i/mochitto/.env

# オーディオデバイスへのアクセス
SupplementaryGroups=audio

# ログ
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mochitto-client

[Install]
WantedBy=multi-user.target
